# configs/etc/profile.d/wb-nmtui.sh
# Alias to use Wiren Board's nmtui wrapper
alias nmtui='nmtui-wb'

