# configs/etc/systemd/system/NetworkManager.service.d/override.conf
[Service]
# https://gitlab.freedesktop.org/NetworkManager/NetworkManager/-/issues/859
CapabilityBoundingSet=CAP_CHOWN
ExecStartPre=/usr/lib/wb-configs/fix_nm_ap_ssid.sh
