# configs/lib/systemd/system/wb-connection-manager.service.d/override.conf
[Service]
ExecCondition=false  # disable wb-connection-manager in this build
