# configs/etc/systemd/system/hostapd.service.d/override.conf
[Service]
ExecCondition=/bin/bash /usr/lib/wb-configs/hostapd_execcondition.sh
