# configs/etc/systemd/system/hostapd.service.d/override.conf
[Service]
ExecCondition=/usr/lib/wb-configs/hostapd_execcondition.sh
