# configs/etc/systemd/system/ModemManager.service.d/override.conf
[Unit]
# MM should be started before any wb-gsm call (wb-gsm could stop MM due to internal assumptions)
Before=ifupdown-pre.service
