# configs/etc/systemd/system/mosquitto.service.d/override.conf
[Service]
ExecStartPre=/bin/chown mosquitto:mosquitto /run/mosquitto
ExecStartPre=/bin/chmod g+x /run/mosquitto