# configs/lib/systemd/system/ModemManager.service.d/override.conf
[Service]
ExecCondition=/usr/lib/wb-configs/detect_legacy_modem_config.sh
