# configs/lib/systemd/system/ModemManager.service.d/override.conf
[Service]
ExecCondition=false  # disable ModemManager in this build
