# configs/lib/systemd/system/NetworkManager.service.d/override.conf
[Service]
ExecCondition=false  # disable NetworkManager in this build
