# configs/etc/systemd/system/nginx.service.d/restart.conf
[Unit]
StartLimitInterval=120

[Service]
Restart=on-failure
RestartSec=5
