# configs/etc/systemd/system/nginx.service.d/restart.conf
[Service]
Restart=on-failure
RestartSec=5s
RestartMaxDelaySec=120s
RestartSteps=5
