# configs/etc/systemd/system/nginx.service.d/restart.conf
[Service]
Restart=on-failure
RestartSec=5
StartLimitInterval=120
TimeoutStartSec=30
