# configs/etc/systemd/system/fcgiwrap.service.d/override.conf
[Service]
ExecStart=/usr/sbin/fcgiwrap -f -c 2

