# configs/etc/systemd/system/fcgiwrap.service.d/override.conf
[Service]
Environment="DAEMON_OPTS=-f -c 30"
ExecStartPre=+/usr/lib/wb-configs/fix_fastcgi.sh
ExecStart=
ExecStart=/usr/sbin/fcgiwrap $DAEMON_OPTS
